# build/packages/v3xctrl/usr/lib/systemd/system/v3xctrl-setup-env.service
[Unit]
Description=Setup env for all other services
DefaultDependencies=no
RequiresMountsFor=/data
Requires=etc-v3xctrl.mount
After=local-fs.target etc-v3xctrl.mount
Before=multi-user.target

[Service]
Type=oneshot
ExecStart=/usr/bin/v3xctrl-write-env
RemainAfterExit=true
TimeoutStartSec=10s

[Install]
WantedBy=multi-user.target
