# build/packages/v3xctrl/usr/lib/systemd/system/v3xctrl-control-failure-handler.service
[Unit]
Description=Handle  failure

[Service]
Type=oneshot
ExecStart=/usr/bin/v3xctrl-control-failed
User=root

[Install]
WantedBy=multi-user.target
