# build/packages/v3xctrl/usr/lib/systemd/system/v3xctrl-firstboot.service
[Unit]
Description=Run first boot script if present
After=cloud-init.target boot-firmware.mount
ConditionPathExists=/boot/firmware/firstboot.sh

[Service]
Type=oneshot
ExecStart=/boot/firmware/firstboot.sh
RemainAfterExit=true
StandardOutput=journal+console
StandardError=journal+console
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
