# build/packages/v3xctrl/usr/lib/systemd/system/v3xctrl-firstboot.service
[Unit]
Description=Run first boot script if present
ConditionPathExists=/boot/firmware/firstboot.sh
After=local-fs.target basic.target
Before=getty@tty1.service
DefaultDependencies=no

[Service]
Type=oneshot
ExecStart=/boot/firmware/firstboot.sh
RemainAfterExit=true
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
