# build/packages/v3xctrl/DEBIAN/postinst
#!/bin/bash
set -e

NAME="v3xctrl"
USER="v3xctrl"

# Create user and add to groups
if ! getent passwd $USER > /dev/null; then
  useradd --system --no-create-home --shell /usr/sbin/nologin $USER
fi

usermod -aG render,video,dialout,systemd-journal,i2c,gpio $USER

# Make sure user can write config dir and files (source of bind mount)
chown $USER:$USER "/data/config"
chown $USER:$USER "/data/config/config.json"
chmod 640 "/data/config/config.json"

# Grant capability to bind to privileged ports (port 80)
setcap 'cap_net_bind_service=+ep' /opt/v3xctrl-python/bin/python3.11

# Update env exports
v3xctrl-write-env

if command -v systemctl >/dev/null 2>&1; then
  # Disable services that we do not need or are handled by other services
  systemctl disable \
    rpi-eeprom-update.service \
    rsync.service \
    udisks2.service \
    apt-daily.timer apt-daily-upgrade.timer \
    man-db.timer \
    e2scrub_all.timer e2scrub_reap.service \
    fstrim.timer \
    remote-fs.target \
    nmbd.service \
    keyboard-setup.service \
    bluetooth.service \
    ModemManager.service

  # Enable our services on start
  systemctl enable \
    etc-v3xctrl.mount \
    v3xctrl-config-server.service \
    v3xctrl-service-manager.service \
    v3xctrl-wifi-mode.service \
    v3xctrl-setup-env.service

  # Make sure the minimum requirements to start the chain after install are in
  # place.
  systemctl restart \
    etc-v3xctrl.mount \
    v3xctrl-config-server.service
fi
