# build/packages/v3xctrl/DEBIAN/prerm
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  # Reset services only on removal
  if [ "$1" = "remove" ]; then
    systemctl disable \
      etc-v3xctrl.mount \
      v3xctrl-config-server.service \
      v3xctrl-service-manager.service \
      v3xctrl-wifi-mode.service \
      v3xctrl-setup-env.service

    # Re-enable services we disabled during install
    systemctl enable \
      rpi-eeprom-update.service \
      rsync.service \
      udisks2.service \
      apt-daily.timer apt-daily-upgrade.timer \
      man-db.timer \
      e2scrub_all.timer e2scrub_reap.service \
      fstrim.timer \
      remote-fs.target \
      nmbd.service \
      keyboard-setup.service \
      bluetooth.service \
      ModemManager.service
  fi
fi
