# build/v3xctrl/DEBIAN/prerm
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  # Reset services only on removal
  if [ "$1" = "remove" ]; then
    systemctl disable \
      rc-config-server.service \
      rc-service-manager.service \
      rc-cleanup.service \
      v3xctrl-wifi-mode.service \
      v3xctrl-setup-env.service

    # Re-enable services we disabled during install
    systemctl enable \
      wpa_supplicant.service \
      hostapd.service \
      dnsmasq.service \
      rpi-eeprom-update.service \
      rsync.service \
      triggerhappy.service triggerhappy.socket \
      udisks2.service \
      apt-daily.timer apt-daily-upgrade.timer \
      man-db.timer \
      e2scrub_all.timer e2scrub_reap.service \
      fstrim.timer \
      rsyslog.service \
      rpi-display-backlight.service \
      remote-fs.target \
      nfs-client.target \
      nmbd.service \
      keyboard-setup.service \
      rng-tools-debian.service \
      bluetooth.service \
      bthelper@hci0.service \
      hciuart.service \
      ModemManager.service
  fi
fi
