# build/v3xctrl/lib/systemd/system/rc-control-failure-handler.service
[Unit]
Description=Handle rc-control failure

[Service]
Type=oneshot
ExecStart=/usr/bin/rc-control-failed
User=root

[Install]
WantedBy=multi-user.target
