# build/v3xctrl/lib/systemd/system/rc-config-server.service
[Unit]
Description=RC Config Server
After=network.target rc-wifi-mode.service
Requires=rc-wifi-mode.service

[Service]
Environment=PYENV_ROOT=/usr/local/pyenv
Environment=PATH=/usr/local/pyenv/shims:/usr/local/pyenv/bin:/usr/bin:/bin
ExecStart=/usr/bin/rc-config-server
Restart=on-failure
User=rc-user
Group=rc-user
KillSignal=SIGINT
WorkingDirectory=/usr/bin
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
