# build/v3xctrl/lib/systemd/system/rc-control.service
[Unit]
Description=Control channel
PartOf=rc-service-manager.service
After=network.target

[Service]
ExecStart=/usr/bin/rc-control
Restart=on-failure
Type=simple
User=v3xctrl
Group=v3xctrl

# Real-time priority
CPUSchedulingPolicy=rr
CPUSchedulingPriority=55

# Highest I/O priority
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Lock memory to avoid swapping
LimitMEMLOCK=infinity

# Max allowed real-time priority
LimitRTPRIO=infinity
