# build/v3xctrl/lib/systemd/system/rc-service-manager.service
[Unit]
Description=Service Manager (autostarts video and control services)
After=network.target v3xctrl-wifi-mode.service v3xctrl-setup-env.service
Requires=v3xctrl-wifi-mode.service pigpiod.service v3xctrl-setup-env.service

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/rc-service-manager
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
