# build/v3xctrl/lib/systemd/system/rc-wifi-mode.service
[Unit]
Description=Switch between Wi-Fi client and AP mode
After=local-fs.target network-pre.target
Before=network.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/rc-wifi-mode
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
