# build/v3xctrl/lib/systemd/system/v3xctrl-pigpiod.service
[Unit]
Description=Daemon for controlling GPIO pins via pigpio
After=network.target

[Service]
ExecStart=/usr/bin/pigpiod -l
Type=forking
PIDFile=/run/pigpio.pid
KillSignal=SIGKILL
KillMode=process

[Install]
WantedBy=multi-user.target
