# build/v3xctrl/lib/systemd/system/v3xctrl-control.service
[Unit]
Description=Control channel
PartOf=rc-service-manager.service
OnFailure=rc-control-failure-handler.service

[Service]
ExecStart=/usr/bin/v3xctrl-control
Restart=no
Type=simple
User=v3xctrl
Group=v3xctrl

# Real-time priority
CPUSchedulingPolicy=rr
CPUSchedulingPriority=55

# Highest I/O priority
IOSchedulingClass=realtime
IOSchedulingPriority=0

# Lock memory to avoid swapping
LimitMEMLOCK=infinity

# Max allowed real-time priority
LimitRTPRIO=infinity
