# build/packages/v3xctrl/etc/NetworkManager/dispatcher.d/50-no-default-modem-route
#!/bin/bash

# Prevent non-wlan interfaces (e.g. RNDIS modem on eth0) from taking the
# default route or DNS priority unless the user has explicitly configured
# rndis routing. Without this, NetworkManager gives eth0 a lower metric than
# wlan0, causing all traffic and DNS to go through a modem that may have no
# internet (no SIM).

INTERFACE="$1"
ACTION="$2"

if [[ "$INTERFACE" == "wlan0" || "$INTERFACE" == "lo" ]]; then
    exit 0
fi

if [[ "$ACTION" == "up" || "$ACTION" == "dhcp4-change" ]]; then
    ROUTING_MODE="wlan"
    if [ -f /run/v3xctrl.env ]; then
        ROUTING_MODE=$(grep -oP '^network_routing=\K.*' /run/v3xctrl.env || echo "wlan")
    fi

    if [[ "$ROUTING_MODE" != "rndis" ]]; then
        ip route del default dev "$INTERFACE" 2>/dev/null || true

        # Remove DNS servers provided by this interface so wlan0 DNS is used.
        # NM regenerates resolv.conf from remaining connections.
        CONNECTION=$(nmcli -g GENERAL.CONNECTION device show "$INTERFACE" 2>/dev/null)
        if [ -n "$CONNECTION" ]; then
            nmcli device modify "$INTERFACE" ipv4.dns "" 2>/dev/null || true
            nmcli device modify "$INTERFACE" ipv6.dns "" 2>/dev/null || true
        fi
    fi
fi
