# build/rc-client/DEBIAN/postinst
#!/bin/bash
set -e

USER=rc-user

# Create user and add to groups
if ! getent passwd $USER > /dev/null; then
  useradd --system --no-create-home --shell /usr/sbin/nologin $USER
fi
usermod -aG video $USER

chown $USER:$USER /etc/rc-client
find /etc/rc-client -maxdepth 1 -name '*.json' -exec chown $USER:$USER {} +
chmod 640 /etc/rc-client/config.json

RECORD_DIR="/var/lib/rc-client/recordings"
if [ ! -d "$RECORD_DIR" ]; then
    mkdir -p "$RECORD_DIR"
    chown $USER:$USER "$RECORD_DIR"
    chmod a+rw "$RECORD_DIR"
fi

if command -v systemctl >/dev/null 2>&1; then
  deb-systemd-helper daemon-reload

  # Disable services that we do not need or are handled by other services
  deb-systemd-helper disable \
    dhcpcd \
    wpa_supplicant \
    hostapd \
    dnsmasq \
    bluetooth \
    rpi-eeprom-update \
    rsync \
    triggerhappy triggerhappy.socket \
    udisks2 \
    apt-daily.timer apt-daily-upgrade.timer \
    man-db.timer \
    e2scrub_all.timer e2scrub_reap \
    fstrim.timer \
    rsyslog \
    rpi-display-backlight \
    remote-fs.target \
    nfs-client.target

  # Enable our services on start
  deb-systemd-helper enable \
    rc-config-server \
    rc-service-manager \
    rc-wifi-mode \
    pigpiod

  # Make sure the minimum requirements to sart the chain after install are in
  # place.
  deb-systemd-helper restart rc-config-server
  deb-systemd-helper restart pigpiod
fi
