# build/rc-client/DEBIAN/postinst
#!/bin/bash

if ! id -u rc-user >/dev/null 2>&1; then
  useradd --system --no-create-home --shell /usr/sbin/nologin rc-user
  sudo usermod -aG video rc-user
fi

chown rc-user:rc-user /etc/rc-client
chown rc-user:rc-user /etc/rc-client/*.json
chmod a+r /etc/rc-client/config.json

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload

  systemctl enable rc-config-server
  systemctl enable rc-service-manager
  systemctl enable rc-wifi-mode

  # Services that are handled by one of our services
  systemctl disable dhcpcd
  systemctl disable wpa_supplicant
  systemctl disable hostapd
  systemctl disable dnsmasq

  # Disable services that we don not need for our usecase
  systemctl disable \
    bluetooth \
    rpi-eeprom-update \
    rsync \
    triggerhappy triggerhappy.socket \
    udisks2 \
    apt-daily.timer apt-daily-upgrade.timer \
    man-db.timer \
    e2scrub_all.timer e2scrub_reap \
    fstrim.timer \
    rsyslog \
    rpi-display-backlight \
    remote-fs.target \
    nfs-client.target
fi
