# build/rc-client/DEBIAN/postinst
#!/bin/bash

if ! id -u rc-user >/dev/null 2>&1; then
  useradd --system --no-create-home --shell /usr/sbin/nologin rc-user
fi

chown rc-user:rc-user /etc/rc-client
chown rc-user:rc-user /etc/rc-client/*.json
chmod a+r /etc/rc-client/config.json

if command -v systemctl >/dev/null 2>&1; then
  systemctl daemon-reload

  systemctl enable rc-config-server
  systemctl start rc-config-server

  systemctl enable rc-service-manager
  systemctl enable rc-wifi-mode
fi
