# build/rc-client/DEBIAN/postinst
#!/bin/bash
set -e

USER=rc-user

# Create user and add to groups
if ! getent passwd $USER > /dev/null; then
  useradd --system --no-create-home --shell /usr/sbin/nologin $USER
fi
usermod -aG video $USER

chown $USER:$USER /etc/rc-client
find /etc/rc-client -maxdepth 1 -name '*.json' -exec chown $USER:$USER {} +
chmod 640 /etc/rc-client/config.json

RECORD_DIR="/var/lib/rc-client/recordings"
if [ ! -d "$RECORD_DIR" ]; then
    mkdir -p "$RECORD_DIR"
    chown $USER:$USER "$RECORD_DIR"
    chmod a+rw "$RECORD_DIR"
fi

if command -v systemctl >/dev/null 2>&1; then
  # Disable services that we do not need or are handled by other services
  deb-systemd-helper disable \
    dhcpcd.service \
    wpa_supplicant.service \
    hostapd.service \
    dnsmasq.service \
    bluetooth.service \
    rpi-eeprom-update.service \
    rsync.service \
    triggerhappy.service triggerhappy.socket \
    udisks2.service \
    apt-daily.timer apt-daily-upgrade.timer \
    man-db.timer \
    e2scrub_all.timer e2scrub_reap.service \
    fstrim.timer \
    rsyslog.service \
    rpi-display-backlight.service \
    remote-fs.target \
    nfs-client.target

  # Enable our services on start
  deb-systemd-helper enable \
    rc-config-server.service \
    rc-service-manager.service \
    rc-wifi-mode.service \
    pigpiod.service

  # Make sure the minimum requirements to start the chain after install are in
  # place.
  deb-systemd-invoke restart \
    rc-config-server.service \
    pigpiod.service
fi
