# build/rc-client/DEBIAN/prerm
#!/bin/bash
systemctl stop rc-config-server
systemctl disable rc-config-server
