# build/rc-client/DEBIAN/prerm
#!/bin/bash
set -e

if command -v systemctl >/dev/null 2>&1; then
  # Reset services only on removal
  if [ "$1" = "remove" ]; then
    deb-systemd-helper disable \
      rc-config-server \
      rc-service-manager \
      rc-wifi-mode

    # Re-enable services we disabled during install
    deb-systemd-helper enable \
      dhcpcd \
      wpa_supplicant \
      hostapd \
      dnsmasq \
      bluetooth \
      rpi-eeprom-update \
      rsync \
      triggerhappy triggerhappy.socket \
      udisks2 \
      apt-daily.timer apt-daily-upgrade.timer \
      man-db.timer \
      e2scrub_all.timer e2scrub_reap \
      fstrim.timer \
      rsyslog \
      rpi-display-backlight \
      remote-fs.target \
      nfs-client.target
  fi
fi
