# build/rc-client/DEBIAN/prerm
#!/bin/bash

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop rc-config-server
  systemctl disable rc-config-server
fi
