# build/packages/v3xctrl/lib/systemd/system/v3xctrl-wifi-mode.service
[Unit]
Description=Switch between Wi-Fi client and AP mode
After=sys-subsystem-net-devices-wlan0.device
Requires=sys-subsystem-net-devices-wlan0.device

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/v3xctrl-wifi-mode
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
