# build/packages/v3xctrl/lib/systemd/system/v3xctrl-cleanup.service
[Unit]
Description=Recover from unclean shutdown
DefaultDependencies=no
Before=sysinit.target
Wants=network-pre.target

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/v3xctrl-cleanup
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
