# build/rc-client/lib/systemd/system/rc-control.service
[Unit]
Description=RC Control channel
After=network.target

[Service]
ExecStart=/usr/bin/rc-control
Restart=on-failure
User=rc-user
Group=rc-user

[Install]
WantedBy=multi-user.target
