# build/rc-client/lib/systemd/system/rc-wifi-mode.service
[Unit]
Description=Switch between Wi-Fi client and AP mode
After=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/bin/rc-wifi-mode
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
